# conf/supervisor/fastlab.conf
[program:fastlab]
command=php -S 127.0.0.1:8080 -t [dir]
autostart=true
autorestart=true
stderr_logfile=/var/log/fastlab.err.log
stdout_logfile=/var/log/fastlab.out.log
user=root